# config/profile.xml
<?xml version='1.0'?>
<!--
    Copyright 2023 Oxide Computer Company
-->
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='profile' name='platform_omicron1'>

<!--
    Make sure IPv4 routing is disabled.  (See also: illumos bug 14006.)
-->
    <service name='network/routing/route' version='1' type='service'>
        <instance name='default' enabled='false' />
    </service>
    <service name='network/routing-setup' version='1' type='service'>
        <instance name='default' enabled='true'>
            <property_group name='routeadm' type='application'>
                <propval name='ipv4-routing-set' type='boolean' value='true' />
            </property_group>
        </instance>
    </service>

<!--
    Disable the dns/client service by default. Most Omicron zones do not need
    this and those that do will configure and then enable it as part of
    service setup:
-->
    <service name='network/dns/client' version='1' type='service'>
        <instance name='default' enabled='false' />
    </service>

<!--
    Disable any other services that we do not expect to need within an Omicron
    zone:
-->
    <service name='network/ipmp' version='1' type='service'>
        <instance name='default' enabled='false' />
    </service>
    <service name='network/rpc/bind' version='1' type='service'>
        <instance name='default' enabled='false' />
    </service>
    <service name='network/rpc/gss' version='1' type='service'>
        <instance name='default' enabled='false' />
    </service>
    <service name='network/security/ktkt_warn' version='1' type='service'>
        <instance name='default' enabled='false' />
    </service>
    <service name='system/sac' version='1' type='service'>
        <instance name='default' enabled='false' />
    </service>
    <service name='network/inetd' version='1' type='service'>
        <instance name='default' enabled='false' />
    </service>
</service_bundle>
